# overlay-adv/etc/xdg/autostart/udiskie-automount.desktop
[Desktop Entry]
Name=automounter
Comment=udisks2 automounter with optional notifications;tray icon
Exec=udiskie --tray
